# influxdb/rootfs/etc/services.d/chronograf/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2191
# ==============================================================================
# Home Assistant Community Add-on: InfluxDB
# Runs the Chronograf Server
# ==============================================================================
declare -a options
declare ingress_entry
declare secret

# Wait for InfluxDB to become available
s6-svwait -u -t 5000 /var/run/s6/services/influxdb

# Wait for Kapacitor to become available
s6-svwait -u -t 5000 /var/run/s6/services/kapacitor

secret=$(</data/secret)

# Wait for InfluxDB...
bashio::log.info "Chronograf is waiting until InfluxDB is available..."
for i in {1800..0}; do
    if influx \
        -execute "SHOW DATABASES" \
        -username "chronograf" \
        -password "${secret}" \
        > /dev/null 2>&1; then
        break;
    fi
    sleep 5
done

if [[ "$i" = 0 ]]; then
    bashio::log.warning "Chronograf is done waiting. Starting anyway..."
fi

bashio::log.info 'Starting Chronograf...'

options+=(--host="127.0.0.1")
options+=(--port=8889)
options+=(--bolt-path=/data/chronograf.db)
options+=(--influxdb-url="http://localhost:8086")

options+=(--influxdb-username="chronograf")
options+=(--influxdb-password="${secret}")
options+=(--kapacitor-url="http://localhost:9092")

# Find the matching Chronograph log level
if bashio::config.has_value 'log_level'; then
    case "$(bashio::string.lower "$(bashio::config 'log_level')")" in
        all|trace|debug)
            options+=(--log-level=debug)
            ;;
        info|notice)
            options+=(--log-level=info)
            ;;
        warning|error|fatal|off)
            options+=(--log-level=error)
            ;;
    esac
fi

ingress_entry=$(bashio::addon.ingress_entry)
options+=(--basepath "${ingress_entry}")

# Run Chronograf
exec chronograf "${options[@]}"
