# influxdb/rootfs/etc/services.d/influxdb/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: InfluxDB
# Runs the InfluxDB Server
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

hass.log.info 'Starting the InfluxDB'

# Run InfluxDB
exec influxd
