# influxdb/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: InfluxDB
# Runs the Nginx daemon
# ==============================================================================
# Wait for Chronograf to become available
s6-svwait -u -t 5000 /var/run/s6/services/chronograf

bashio::log.info "Starting NGinx server..."

exec nginx -g "daemon off;"
