# influxdb/rootfs/etc/services.d/kapacitor/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: InfluxDB
# Runs the Kapacitor Server
# ==============================================================================
declare secret

# Wait for InfluxDB to become available
s6-svwait -u -t 5000 /var/run/s6/services/influxdb

secret=$(</data/secret)

# Wait for InfluxDB...
bashio::log.info "Kapacitor is waiting until InfluxDB is available..."
for i in {1800..0}; do
    if influx \
        -execute "SHOW DATABASES" \
        -username "kapacitor" \
        -password "${secret}" \
        > /dev/null 2>&1; then
        break;
    fi
    sleep 5
done

if [[ "$i" = 0 ]]; then
    bashio::log.warning "Kapacitor is done waiting. Starting anyway..."
fi

bashio::log.info 'Starting the Kapacitor'

# Run Kapacitor
exec kapacitord
